# software/configs/ctrl.service
[Unit]
Description=The AXIOM ctrl daemon to expose a fuse api

[Service]
ExecStart=/usr/axiom/bin/ctrl --mountpoint /axiom-api/ $(realpath /etc/axiom-yml)
StandardOutput=journal

[Install]
WantedBy=multi-user.target
